# conf/prism/prism_back.conf

start on runlevel [2345]
stop on runlevel [!2345]
respawn

exec su - fold -c "cd /opt/prism; FOLD_ENV=production bundle exec bin/prism_back 2>&1 | logger -t 'pro|prism_back'"